# backend/deployment/systemd/postgresql.service
[Unit]
Description=PostgreSQL Database Server with pgvector
After=network.target

[Service]
Type=notify
User=postgres
Group=postgres

# Location of database directory
Environment=PGDATA=/var/lib/postgresql/data

# Where to send early-startup messages from the server
# This is normally controlled by the global default set by systemd
StandardOutput=journal

ExecStart=/usr/lib/postgresql/15/bin/postgres -D ${PGDATA}
ExecReload=/bin/kill -HUP $MAINPID
KillMode=mixed
KillSignal=SIGINT
TimeoutSec=infinity

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/var/lib/postgresql /var/log/postgresql

Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target

